# crontab
0 17 * * * /usr/bin/timeshift --create --comments "Daily $(date +'%F')" --tags D
